# modules/profile/files/github_runner/gha-restart-when-idle.service
[Unit]
Description=Restart the GitHub Actions runner to load patched libraries, only when idle

[Service]
Type=oneshot
ExecStart=/usr/local/bin/gha-restart-when-idle.sh
