# environments/sandbox/modules/profile/templates/github_runner/actions-runner.service.erb
[Unit]
Description=GitHub self-hosted runner
After=network.target

[Service]
Type=simple
ExecStart=<%= @runner_package_directory %>/run.sh
WorkingDirectory=<%= @runner_package_directory %>
User=<%= @github_runner_user %>
Group=<%= @github_runner_group %>
Restart=on-failure

[Install]
WantedBy=multi-user.target
