# environments/development/modules/profile/templates/github_runner/actions-runner.service.erb
[Unit]
Description=GitHub self-hosted runner
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=<%= @start_script %>
ExecStopPost=/usr/local/bin/gha-on-runner-exit.sh
Environment=DEREGISTRATION_HOOK_NAME=<%= @deregistration_hookname %>
WorkingDirectory=<%= @runner_package_directory %>
User=<%= @github_runner_user %>
Group=<%= @github_runner_group %>
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=21600
Restart=on-failure

[Install]
WantedBy=multi-user.target
